# machine_code/urbanketl-machine.service
[Unit]
Description=UrbanKetl Tea Machine Controller
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/urbanketl-machine
ExecStart=/usr/bin/python3 /home/pi/urbanketl-machine/urbanketl_machine.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/urbanketl-machine

[Install]
WantedBy=multi-user.target